# daemon/com.user.claude-usage-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.claude-usage-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON_BIN__</string>
        <string>__DAEMON_PATH__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__LOG_OUT__</string>

    <key>StandardErrorPath</key>
    <string>__LOG_ERR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>DEEPSEEK_API_KEY</key>
        <string>__DEEPSEEK_API_KEY__</string>
        <key>OPENAI_API_KEY</key>
        <string>__OPENAI_API_KEY__</string>
        <key>TRANSCRIBE_PROVIDER</key>
        <string>__TRANSCRIBE_PROVIDER__</string>
        <key>DASHSCOPE_API_HOST</key>
        <string>__DASHSCOPE_API_HOST__</string>
        <key>DASHSCOPE_API_KEY</key>
        <string>__DASHSCOPE_API_KEY__</string>
        <key>DASHSCOPE_COMPATIBLE_BASE_URL</key>
        <string>__DASHSCOPE_COMPATIBLE_BASE_URL__</string>
        <key>DASHSCOPE_TRANSCRIBE_MODE</key>
        <string>__DASHSCOPE_TRANSCRIBE_MODE__</string>
        <key>DASHSCOPE_TRANSCRIBE_MODEL</key>
        <string>__DASHSCOPE_TRANSCRIBE_MODEL__</string>
        <key>DASHSCOPE_ASR_LANGUAGE</key>
        <string>__DASHSCOPE_ASR_LANGUAGE__</string>
        <key>DASHSCOPE_TRANSCRIBE_TIMEOUT</key>
        <string>__DASHSCOPE_TRANSCRIBE_TIMEOUT__</string>
        <key>DASHSCOPE_TTS_MODEL</key>
        <string>__DASHSCOPE_TTS_MODEL__</string>
        <key>DASHSCOPE_TTS_VOICE</key>
        <string>__DASHSCOPE_TTS_VOICE__</string>
    </dict>
</dict>
</plist>
